# supervisor/swiftearn-webhook.conf
[program:swiftearn-webhook]
command=node /var/www/html/swiftearn-app/webhook-server.js
directory=/var/www/html/swiftearn-app
user=www-data
autostart=true
autorestart=true
startretries=5
environment=NODE_ENV="production",WEBHOOK_SECRET="YOUR_SECRET_HERE",WEBHOOK_PORT="9001",WEBHOOK_BRANCH="refs/heads/main",WEBHOOK_SCRIPT="/var/www/html/swiftearn-app/deploy-frontend.sh",WEBHOOK_LOG="/var/www/html/swiftearn-app/logs/deploy-frontend.log"
stdout_logfile=/var/log/supervisor/swiftearn-webhook.log
stderr_logfile=/var/log/supervisor/swiftearn-webhook-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
